# roles/networking/avahi/templates/supervisor_conf.d_avahi-daemon.conf.j2
[program:avahi-daemon]
command=/usr/sbin/avahi-daemon --no-rlimits
environment=USER=root,
user=root
