# roles/database/redis/templates/supervisor_conf.d_redis-server.conf.j2
[program:redis-server]
command=/usr/bin/redis-server /etc/redis/redis.conf
environment=USER={{ redis_user }},HOME={{ redis_path }},
user={{ redis_user }}
