# roles/web/nginx/templates/supervisor_conf.d_nginx.conf.j2
[program:nginx]
command=/usr/sbin/nginx
directory={{ nginx_path }}
environment=USER=root,HOME={{ nginx_path }}
user=root
