# roles/web/php5/templates/supervisor_conf.d_php5-fpm.conf.j2
[program:php5-fpm]
command=/usr/sbin/php5-fpm --nodaemonize --fpm-config /etc/php5/fpm/php-fpm.conf
environment=USER={{ php5_user }},HOME={{ php5_path }}
user={{ php5_user }}
