# roles/database/mongodb/templates/supervisor_conf.d_mongod.conf.j2
[program:mongod]
command=/usr/bin/mongod --config /etc/{{ mongodb_service }}.conf
environment=USER={{ mongodb_user }},HOME={{ mongodb_path }},
user={{ mongodb_user }}
